# service/openrc/init.d/tty-colorscheme
#!/sbin/openrc-run

supervisor="supervise-daemon"
description="Change TTY colors"
command="/usr/bin/tty-colorscheme"
command_args="${colorscheme}"

start_pre() {
	${command} ${command_args} &> /dev/null
	eend $? "$(${command} ${command_args})"
}

start() {
	if [ -z "${command_args}" ]; then
		ewarn "No colorscheme set, using system default colors"
	else
		einfo "$(TERM= ${command} ${command_args})"
	fi
}
